# scripts/systemd/testdashboard.service
[Unit]
Description=Replay Test Dashboard (Docker Compose)
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
ExecStart=/bin/bash /opt/replay/scripts/deploy.sh
ExecStop=/usr/bin/docker compose -f /opt/replay/docker/docker-compose.yml down
RemainAfterExit=yes
WorkingDirectory=/opt/replay
User=root

[Install]
WantedBy=multi-user.target
