# scripts/systemd/testdashboard.service
[Unit]
Description=PHIRO Playwright Test Dashboard
After=network.target docker.service
Requires=docker.service

[Service]
Type=simple
User=xbata-mm84
WorkingDirectory=/home/xbata-mm84/testdashboard/docker
EnvironmentFile=/home/xbata-mm84/testdashboard/.env
ExecStartPre=-/usr/bin/docker compose down
ExecStart=/usr/bin/docker compose up
ExecStop=/usr/bin/docker compose down
Restart=always
RestartSec=10s
StartLimitBurst=5
StartLimitIntervalSec=120s
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
